# .config/systemd/user/update-reminder.timer
[Unit]
Description=Run Weekly Arch Linux Update Reminder

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

