# .config/systemd/user/update-reminder.service
[Unit]
Description=Weekly Arch Linux Update Reminder

[Service]
Type=oneshot
ExecStart=%h/.local/bin/update-reminder.sh

